# FreeBSD/files/usr/local/etc/periodic/daily/aide
#!/bin/sh

aide -c $_CONF_INTEGRITY_CONFIGURATION_FILE --check
